# infra/systemd/proxima-backup.service.example
[Unit]
Description=Proxima database backup

[Service]
Type=oneshot
User=proxima
Group=proxima
EnvironmentFile=/etc/proxima/proxima.env
WorkingDirectory=/opt/proxima
ExecStart=/opt/proxima/scripts/backup
