# infra/systemd/proxima.service.example
[Unit]
Description=Proxima self-hosted workspace
After=network.target

[Service]
Type=simple
User=proxima
Group=proxima
EnvironmentFile=/etc/proxima/proxima.env
WorkingDirectory=/opt/proxima
ExecStart=/usr/local/bin/proxima serve
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
